# install/autostart/obj_det_autostart.service
[Unit]
Description=Autostart object detection service
After=network.target
[Service]
ExecStart=/home/orangepi/miniconda3/envs/rknn/bin/python3 /home/orangepi/workspace/Yolov5_RK3588/main.py --webui
Restart=always
[Install]
WantedBy=multi-user.target