# install/autostart/obj_det_autostart.service
[Unit]
Description=Autostart object detection service
After=network.target
[Service]
ExecStart=
Restart=always
[Install]
WantedBy=multi-user.target